# supervisord.conf
[supervisord]
nodaemon=false
logfile=/home/user/webapp/hr-recruitment-dashboard/supervisord.log
pidfile=/home/user/webapp/hr-recruitment-dashboard/supervisord.pid
childlogdir=/home/user/webapp/hr-recruitment-dashboard/logs

[unix_http_server]
file=/home/user/webapp/hr-recruitment-dashboard/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///home/user/webapp/hr-recruitment-dashboard/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:laravel-webserver]
command=php artisan serve --host=0.0.0.0 --port=8000
directory=/home/user/webapp/hr-recruitment-dashboard
autostart=true
autorestart=true
startsecs=3
startretries=3
user=user
redirect_stderr=true
stdout_logfile=/home/user/webapp/hr-recruitment-dashboard/logs/laravel-webserver.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/user/webapp/hr-recruitment-dashboard/logs/laravel-webserver-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

[program:laravel-queue]
command=php artisan queue:work --sleep=3 --tries=3 --timeout=90
directory=/home/user/webapp/hr-recruitment-dashboard
autostart=true
autorestart=true
startsecs=3
startretries=3
user=user
numprocs=1
redirect_stderr=true
stdout_logfile=/home/user/webapp/hr-recruitment-dashboard/logs/laravel-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/user/webapp/hr-recruitment-dashboard/logs/laravel-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5